# configs/kmonad/kmonad_sculpt.service
[Unit]
Description=kmonad keyboard config

[Service]
ExecStart=/usr/bin/kmonad %E/kmonad/sculpt.cfg
Nice=-20

[Install]
WantedBy=default.target
